# templates/default/certbot/certbot.cron.erb
37 3 * * * root test -e /usr/bin/certbot && /usr/bin/chpst -e /etc/container_environment /usr/bin/certbot renew --post-hook "/usr/bin/sv restart opensmtpd"
